# hass-ci-shutdown.timer
[Unit]
Description=HomeAssistant CI shutdown server

[Timer]
OnBootSec=30m

[Install]
WantedBy=timers.target
